# deploy/com.viton.daa-yield-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.viton.daa-yield-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/viton/Data-Analysis-Agent/.venv/bin/gunicorn</string>
        <string>--config</string>
        <string>/Users/viton/Data-Analysis-Agent/gunicorn_conf.py</string>
        <string>api:create_app()</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/viton/Data-Analysis-Agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/viton/.hermes/logs/daa_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/viton/.hermes/logs/daa_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>AGENT_PORT</key>
        <string>5001</string>
    </dict>
</dict>
</plist>
